# server-configs/supervisor-app.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0

autostart=true
autorestart=true
killasgroup=true
stopasgroup=true


[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /usr/local/etc/uwsgi/uwsgi-app.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
